# franklin.service
[Unit]
Description=Franklin Agent
After=network.target

[Service]
Type=simple
Environment="PATH=/home/mjscully/.local/bin:/home/mjscully/.npm-global/bin:/home/mjscully/.bun/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/mjscully/franklin/.env
WorkingDirectory=/home/mjscully/franklin
ExecStart=/usr/bin/npm run franklin
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
